# login-watch/login-watch.service
[Unit]
Description=hartguard login-watch — SSH brute force detector
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/login-watch --config /etc/hartguard/login-watch.yaml
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
