# host-watch/host-watch.service
[Unit]
Description=hartguard host-watch — disk, memory, container, and service monitor
After=network.target docker.service

[Service]
Type=simple
ExecStart=/usr/local/bin/host-watch --config /etc/hartguard/host-watch.yaml
Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
